# package/DEBIAN/postinst
#!/bin/sh

systemctl enable autorun.service
systemctl start autorun.service


