# package/etc/systemd/system/autorun.service
[Unit]
Description=liteSh_server

[Service]
ExecStart=/usr/local/bin/litesh -s
restart=always
[Install]
WantedBy=default.target
